# root/etc/services.d/rtsp2mjpg/run
#!/usr/bin/with-contenv bash
cd /rtsp2mjpg
/rtsp2mjpg/simple-cgi-server